# udev.rules
KERNEL=="sr0", ACTION=="change", RUN+="/bin/bash /usr/bin/riplauncher '%N'"
